# deployment/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis.err
priority=10

[program:blacklist]
command=python main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/blacklist.log
stderr_logfile=/app/logs/blacklist.err
environment=PYTHONUNBUFFERED="1",FLASK_ENV="production",PORT="2541",REDIS_URL="redis://localhost:6379/0"
priority=20
startsecs=10

[group:blacklist-services]
programs=redis,blacklist
priority=999