# offline-packages/blacklist-offline-package-v2.0.0/source-code/scripts/blacklist-deployment-worker.service
[Unit]
Description=Blacklist Deployment Queue Worker
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=jclee
Group=jclee
WorkingDirectory=/home/jclee/app/blacklist
ExecStart=/home/jclee/app/blacklist/scripts/deployment-buffer.sh worker
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blacklist-deployment

# 환경 변수
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="KUBECONFIG=/home/jclee/.kube/config"

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true

# 리소스 제한
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target